# scripts/stash-watcher.service
[Unit]
Description=Stash media file watcher
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=brandon
ExecStart=/home/brandon/projects/docker/scripts/stash-watcher.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
